# roles/personal-setup-linux/templates/login.service
[Unit]
Description = Login Service

[Service]
StandardError=journal
User={{ansible_user}}
ExecStart=/bin/bash --login

[Install]
WantedBy=multi-user.target