# roles/environment/systems/linux/variations/common/files/login.service
[Unit]
Description = Login Service

[Service]
StandardError=journal
User={{username}}
ExecStart=/bin/bash --login

[Install]
WantedBy=multi-user.target